# deploy/systemd/cfm-api.service.in
[Unit]
Description=CFM API Service
After=network.target

[Service]
Type=simple
User={{ CFM_USER }}
Group={{ CFM_USER }}
WorkingDirectory={{ REPO_PATH }}
ExecStart={{ VENV_PATH }}/bin/cfm serve -c {{ CONFIG_PATH }} --host {{ LISTEN_HOST }} --port {{ LISTEN_PORT }}
Restart=always
RestartSec=5
StandardOutput=append:{{ LOG_DIR }}/api.log
StandardError=append:{{ LOG_DIR }}/api-error.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
