# deploy/systemd/cfm-sync.timer.in
[Unit]
Description=CFM Sync Timer (daily)

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
